# root/etc/services.d/vnc-server/run
#!/usr/bin/with-contenv /usr/bin/execlineb
/bin/sh -c "/usr/local/bin/x11vnc -forever -shared -rfbport 5900 -display $DISPLAY -passwd $VNC_PASSWORD"
