# root/etc/services.d/xvfb/run
#!/usr/bin/with-contenv /usr/bin/execlineb
/bin/sh -c "/usr/bin/Xvfb $DISPLAY -ac -screen 0 $SCREEN_DIMENSION"
